# deploy/systemd/trading-fast.service
[Unit]
Description=TradingRobot Fast Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/TradingRobot
EnvironmentFile=/home/pi/TradingRobot/deploy/systemd/trading.env
ExecStart=/home/pi/TradingRobot/.venv/bin/python /home/pi/TradingRobot/orchestrator_fast.py --base-dir /home/pi/TradingRobot --config config/config.yaml
Restart=always
RestartSec=15
StandardOutput=append:/home/pi/TradingRobot/logs/fast.service.log
StandardError=append:/home/pi/TradingRobot/logs/fast.service.err.log

[Install]
WantedBy=multi-user.target
